# benchmark_scripts/tpch/tpch_monetdb.service
[Unit]
Description=Run TPCH benchmark for Carbon - MonetDB

[Service]
# Everything in GB
User=michalis
Group=michalis
EnvironmentFile=/etc/systemd/system/tpch_benchmark.env
Environment="DB_FARM=${DATA_DIR}/monetdb_files"

ExecStart= /bin/bash -c 'source ${PYTHON_ENV}/bin/activate && cd ${PROJECT_PATH} && exec ${PYTHON_ENV}/bin/python3 ${PYTHON_PATH} --benchmark tpch --database monetdb --data ${DATA_DIR} --dbname tpch --dbfarm ${DB_FARM}'

#MemoryMax=150G
MemorySwapMax=0
#CPUAffinity=0-13

StandardOutput=file:/home/michalis/Projects/Environmental/Carbon/codebase/olap_db_carbon_benchmark/benchmark_scripts/tpch/logs/tpch_monetdb_output.log
StandardError=file:/home/michalis/Projects/Environmental/Carbon/codebase/olap_db_carbon_benchmark/benchmark_scripts/tpch/logs/tpch_monetdb_error.log
