# benchmark_scripts/tpch/tpch_duckdb.service
[Unit]
Description=Run TPCH benchmark for Carbon - DuckDB
#Not needed anymore
[Service]
# Everything in GB
User=michalis
Group=michalis
EnvironmentFile=/etc/systemd/system/tpch_benchmark.env

ExecStart= /bin/bash -c 'source ${PYTHON_ENV}/bin/activate && cd ${PROJECT_PATH} && exec ${PYTHON_ENV}/bin/python3 ${PYTHON_PATH} --benchmark tpch --database duckdb --data ${DATA_DIR} --db_path ${DATA_DIR}'

#MemoryMax=224G
MemorySwapMax=0
#CPUAffinity=0-13

StandardOutput=file:/home/michalis/Projects/Environmental/Carbon/codebase/olap_db_carbon_benchmark/benchmark_scripts/tpch/logs/tpch_duckdb_output.log
StandardError=file:/home/michalis/Projects/Environmental/Carbon/codebase/olap_db_carbon_benchmark/benchmark_scripts/tpch/logs/tpch_duckdb_error.log
