# benchmark_scripts/tpch/tpch_clickhouse.service
[Unit]
Description=Run TPCH benchmark for Carbon - ClickHouse

[Service]
# Everything in GB
User=michalis
Group=michalis
EnvironmentFile=/etc/systemd/system/tpch_benchmark.env
#Environment="MEM_LIMIT=512"

ExecStart= /bin/bash -c 'source ${PYTHON_ENV}/bin/activate && cd ${PROJECT_PATH} && exec stdbuf -oL ${PYTHON_ENV}/bin/python3 ${PYTHON_PATH} --benchmark tpch --database clickhouse --data ${DATA_DIR}'

#MemoryMax=${MEM_LIMIT}
MemorySwapMax=0
#CPUAffinity=0-13

StandardOutput=file:/home/michalis/Projects/Environmental/Carbon/codebase/olap_db_carbon_benchmark/benchmark_scripts/tpch/logs/tpch_clickhouse_output.log
StandardError=file:/home/michalis/Projects/Environmental/Carbon/codebase/olap_db_carbon_benchmark/benchmark_scripts/tpch/logs/tpch_clickhouse_error.log
